# files/gmnas-welcome.service
[Unit]
Description=gm-nas welcome app (post-setup web UI)
After=network-online.target
Wants=network-online.target
# Only after WiFi setup is done — during first-boot setup, wifi-connect owns
# port 80 for its captive portal. (No conflict / no shadowing.)
ConditionPathExists=/etc/homenas/provisioned

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/gmnas-welcome/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
