# files/homenas-firstboot.service
[Unit]
Description=gm-nas first-boot WiFi provisioning (wifi-connect captive portal)
Documentation=https://github.com/affigabmag/gm-nas
After=NetworkManager.service
Wants=NetworkManager.service
# Skip entirely once the device has been provisioned.
ConditionPathExists=!/etc/homenas/provisioned

[Service]
# Type=simple so it runs in the BACKGROUND: boot continues, the console/SSH
# come up, and wifi-connect broadcasts the setup AP at the same time. (oneshot
# blocked the whole boot until setup finished — the box looked frozen.)
Type=simple
ExecStart=/usr/local/sbin/firstboot-wifi.sh
# NOTE: no Restart= — a failed run (e.g. no AP-capable WiFi adapter) must not
# loop and flood the console. The wrapper handles missing hardware gracefully.

[Install]
WantedBy=multi-user.target
