# files/homenas-firstboot.service
[Unit]
Description=gm-nas first-boot WiFi provisioning (wifi-connect captive portal)
Documentation=https://github.com/affigabmag/gm-nas
After=NetworkManager.service
Wants=NetworkManager.service
# Runs on EVERY boot (no flag gate). The script decides by CONNECTIVITY:
# if the box has an active network it exits immediately (normal boot); if it
# has no working WiFi it launches the setup AP. This guarantees a headless box
# (no keyboard/screen) can ALWAYS be (re)configured from a phone — even if a
# saved network stops working (wrong password, moved home, router replaced).

[Service]
# Type=simple so it runs in the BACKGROUND: boot continues, the console/SSH
# come up, and wifi-connect broadcasts the setup AP at the same time. (oneshot
# blocked the whole boot until setup finished — the box looked frozen.)
Type=simple
ExecStart=/usr/local/sbin/firstboot-wifi.sh
# NOTE: no Restart= — a failed run (e.g. no AP-capable WiFi adapter) must not
# loop and flood the console. The wrapper handles missing hardware gracefully.

[Install]
WantedBy=multi-user.target
